# lesson23_22_06_20/crontab
@hourly /vagrant/mysql_backup_script.sh >/dev/null 2>&1
@hourly /vagrant/pgsql_backup_script.sh >/dev/null 2>&1
