# lesson5_22_04_07/simple-u.service
[Unit]
Description=Test unit service
[Service]
ExecStart=/home/vagrant/share/unit-script.sh
[Install]